# templates/lib/systemd/system/ssh-keygen.target.j2
# {{ ansible_managed }}

[Unit]
Wants=ssh-keygen@dsa.service
Wants=ssh-keygen@ecdsa.service
Wants=ssh-keygen@ed25519.service
Wants=ssh-keygen@rsa.service
PartOf=ssh.service
